# ai-agent.service.template
[Unit]
Description=RAVEN AI Agent - Telegram-controlled local LLM agent
After=network.target

[Service]
Type=simple
User=__WSL_USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=__VENV_DIR__/bin/python3 __PROJECT_DIR__/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
# Clean stop via /stop command or systemctl stop won't trigger restart
# Only unexpected crashes trigger restart (on-failure, not always)
SuccessExitStatus=0
